# deploy/amir-pos.service
# systemd unit for the ERP backend (which also serves the Angular frontend).
# Install:  sudo cp deploy/amir-pos.service /etc/systemd/system/
#           sudo systemctl daemon-reload && sudo systemctl enable --now amir-pos
[Unit]
Description=Amir POS / ClothingERP (Node + Express + Prisma)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/amir-pos/backend
# dotenv reads ./.env from WorkingDirectory, so /opt/amir-pos/backend/.env is loaded.
ExecStart=/usr/bin/node dist/server.js
Restart=always
RestartSec=3
# Give Node a sane heap ceiling on the free ARM box.
Environment=NODE_OPTIONS=--max-old-space-size=1024
StandardOutput=journal
StandardError=journal
SyslogIdentifier=amir-pos

[Install]
WantedBy=multi-user.target
